# backend/slack-memories.service
[Unit]
Description=Slack Memories Backend API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/slack-memories/backend
EnvironmentFile=/opt/slack-memories/backend/.env
ExecStart=/opt/slack-memories/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
